# web/unicorn/gunicorn-restart.service
# /etc/systemd/system/gunicorn-restart.service
[Unit]
Description=Restart gunicorn

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart gunicorn